# mab-approval-daily.timer
[Unit]
Description=Run MAB Approval Daily Tasks

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
