# mab-approval-daily.service
[Unit]
Description=Execute Daily MAB Approval Scripts
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/mab-approval
EnvironmentFile=/opt/mab-approval/.env
ExecStart=/bin/bash /opt/mab-approval/run-daily-tasks.sh
